# packaging/deb/prerm
#!/bin/sh
# cargo-deb maintainer script: stop + disable the service before
# the package files are removed, so an upgrade/removal doesn't leave a running
# unit pointing at a deleted binary.
set -e

if [ "$1" = "remove" ] || [ "$1" = "deconfigure" ]; then
    if [ -d /run/systemd/system ]; then
        systemctl --quiet stop boatramp.service >/dev/null 2>&1 || true
        systemctl --quiet disable boatramp.service >/dev/null 2>&1 || true
    fi
fi

exit 0
